# package/etc/init/mediasmartserver.conf
# upstart script for mediasmartserver

description "Start or stop mediasmartserver"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/sbin/mediasmartserverd --activity --daemon
