# debian/preinst
#!/bin/sh
# preinst script for mediasmartserver
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <new-preinst> `install'
#        * <new-preinst> `install' <old-version>
#        * <new-preinst> `upgrade' <old-version>
#        * <old-preinst> `abort-upgrade' <new-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    install|upgrade)
        if [ -x "/etc/init.d/mediasmartserver" ]
        then
            if [ -x "`which invoke-rc.d 2>/dev/null`" ]
            then
                invoke-rc.d mediasmartserver stop > /dev/null
            else
                /etc/init.d/mediasmartserver stop > /dev/null
            fi

            rm /etc/init.d/mediasmartserver > /dev/null
            update-rc.d mediasmartserver remove > /dev/null
        else
            INIT=$(ps --pid 1 | grep -q systemd && echo 'systemd' || echo 'upstart')
            case "$INIT" in
                systemd)
                    systemctl stop mediasmartserver
                    systemctl disable mediasmartserver
                    systemctl daemon-reload
                    ;;

                *)
                    stop mediasmartserver
                    ;;
            esac
        fi
    ;;

    abort-upgrade)
    ;;

    *)
        echo "preinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
